# launchd/com.USER.usage-hud.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Template. install.sh generates the real file with your username filled in
  and installs it to ~/Library/LaunchAgents/. Replace __USER__ if you install
  by hand.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.__USER__.usage-hud</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/__USER__/.local/bin/usage-hud</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <false/>
  <key>ProcessType</key>
  <string>Interactive</string>
  <key>StandardOutPath</key>
  <string>/tmp/usage-hud.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/usage-hud.err.log</string>
</dict>
</plist>
